# dashboards/persistence-mechanisms.xml
<?xml version="1.0" encoding="UTF-8"?>
<dashboard version="1.1" theme="dark">
  <label>Persistence Mechanisms</label>
  <description>Track persistence techniques including new services, registry modifications, scheduled tasks, and account manipulation. | ATT&amp;CK: T1543.003, T1547.001, T1053.005, T1136.001</description>

  <!-- ============================================================
       Dashboard: Persistence Mechanisms
       Author: Kud Lee
       Created: 2026-03-19
       Log Sources: Windows Security (4720, 4728, 4732, 4756),
                    System (7045), Sysmon (Event 1, 13)
       Purpose: Detect attacker persistence by monitoring service creation,
                registry run keys, scheduled tasks, and user/group changes.
                Critical for identifying post-exploitation activity.
       ============================================================ -->

  <fieldset submitButton="true" autoRun="true">
    <input type="time" token="time_range" searchWhenChanged="true">
      <label>Time Range</label>
      <default>
        <earliest>-24h@h</earliest>
        <latest>now</latest>
      </default>
    </input>
    <input type="dropdown" token="host_filter" searchWhenChanged="true">
      <label>Host</label>
      <choice value="*">All Hosts</choice>
      <default>*</default>
      <search>
        <query>index=wineventlog OR index=sysmon earliest=-24h | stats count by ComputerName | sort - count | head 50</query>
      </search>
      <fieldForLabel>ComputerName</fieldForLabel>
      <fieldForValue>ComputerName</fieldForValue>
    </input>
  </fieldset>

  <!-- ===================== ROW 1: New Services Created (Event 7045) ===================== -->
  <row>
    <panel>
      <title>New Services Created (System Event 7045)</title>
      <table>
        <search>
          <query>index=wineventlog source="WinEventLog:System" EventCode=7045
| eval install_time=strftime(_time, "%Y-%m-%d %H:%M:%S")
| eval service_start_type=case(
    ServiceStartType=0, "Boot",
    ServiceStartType=1, "System",
    ServiceStartType=2, "Automatic",
    ServiceStartType=3, "Manual",
    ServiceStartType=4, "Disabled",
    1=1, ServiceStartType
    )
| eval service_type=case(
    ServiceType="kernel driver", "Kernel Driver",
    ServiceType="user mode service", "User Mode Service",
    1=1, ServiceType
    )
| eval risk=case(
    match(ImagePath, "(?i)(cmd\.exe|powershell\.exe|pwsh\.exe|mshta\.exe|rundll32\.exe)"), "CRITICAL",
    match(ImagePath, "(?i)(temp|appdata|public|downloads|users)"), "HIGH",
    match(ImagePath, "(?i)(\\\\\\\\|http|ftp)"), "HIGH",
    service_start_type="Automatic" AND NOT match(ImagePath, "(?i)(windows|program files|system32)"), "MEDIUM",
    1=1, "LOW"
    )
| table install_time, ComputerName, ServiceName, ImagePath, service_type, service_start_type, AccountName, risk
| rename install_time as "Time", ComputerName as "Host", ServiceName as "Service Name", ImagePath as "Service Binary Path", service_type as "Type", service_start_type as "Start Type", AccountName as "Account", risk as "Risk"
| sort - Time</query>
          <earliest>$time_range.earliest$</earliest>
          <latest>$time_range.latest$</latest>
        </search>
        <option name="count">20</option>
        <option name="drilldown">row</option>
        <option name="wrap">true</option>
        <format type="color" field="Risk">
          <colorPalette type="map">{"CRITICAL":#DC4E41,"HIGH":#F1813F,"MEDIUM":#F8BE34,"LOW":#65A637}</colorPalette>
        </format>
      </table>
    </panel>
  </row>

  <!-- ===================== ROW 2: Registry Run Key Changes (Sysmon 13) ===================== -->
  <row>
    <panel>
      <title>Registry Run Key Changes (Sysmon Event 13)</title>
      <table>
        <search>
          <query>index=sysmon EventCode=13
| where match(TargetObject, "(?i)(\\\\Run\\\\|\\\\RunOnce\\\\|\\\\RunServices\\\\|\\\\Explorer\\\\Shell Folders|\\\\Winlogon\\\\Shell|\\\\Winlogon\\\\Userinit|\\\\CurrentVersion\\\\Image File Execution Options)")
| eval mod_time=strftime(_time, "%Y-%m-%d %H:%M:%S")
| eval process_name=mvindex(split(Image, "\\"), -1)
| eval reg_key=mvindex(split(TargetObject, "\\"), -1)
| eval reg_hive=case(
    match(TargetObject, "HKLM"), "HKLM (Machine)",
    match(TargetObject, "HKU"), "HKU (User)",
    1=1, "Unknown"
    )
| eval risk=case(
    match(Details, "(?i)(cmd\.exe|powershell\.exe|mshta\.exe|rundll32\.exe|wscript\.exe|cscript\.exe)"), "CRITICAL",
    match(Details, "(?i)(temp|appdata|public|downloads)"), "HIGH",
    match(Details, "(?i)(http|\\\\\\\\)"), "HIGH",
    match(TargetObject, "(?i)Image File Execution Options"), "HIGH",
    1=1, "MEDIUM"
    )
| table mod_time, Computer, process_name, EventType, reg_hive, TargetObject, Details, risk
| rename mod_time as "Time", Computer as "Host", process_name as "Modifying Process", EventType as "Action", reg_hive as "Hive", TargetObject as "Registry Key", Details as "Value", risk as "Risk"
| sort - Time</query>
          <earliest>$time_range.earliest$</earliest>
          <latest>$time_range.latest$</latest>
        </search>
        <option name="count">20</option>
        <option name="drilldown">row</option>
        <option name="wrap">true</option>
        <format type="color" field="Risk">
          <colorPalette type="map">{"CRITICAL":#DC4E41,"HIGH":#F1813F,"MEDIUM":#F8BE34}</colorPalette>
        </format>
      </table>
    </panel>
  </row>

  <!-- ===================== ROW 3: Scheduled Tasks Created ===================== -->
  <row>
    <panel>
      <title>Scheduled Tasks Created (Sysmon Event 1 + Security 4698)</title>
      <table>
        <search>
          <query>(index=sysmon EventCode=1 Image="*\\schtasks.exe" CommandLine="*/create*")
    OR (index=wineventlog EventCode=4698)
| eval task_time=strftime(_time, "%Y-%m-%d %H:%M:%S")
| eval source_type=case(
    EventCode=1, "Process Creation (Sysmon)",
    EventCode=4698, "Task Registration (Security)"
    )
| eval task_detail=coalesce(TaskName, CommandLine)
| eval process_name=if(isnotnull(Image), mvindex(split(Image, "\\"), -1), "N/A")
| eval user=coalesce(User, SubjectUserName)
| eval risk=case(
    match(task_detail, "(?i)(powershell|cmd\.exe|mshta|wscript|cscript|rundll32)"), "HIGH",
    match(task_detail, "(?i)(temp|appdata|public|downloads)"), "HIGH",
    match(task_detail, "(?i)(http|\\\\\\\\)"), "HIGH",
    1=1, "MEDIUM"
    )
| table task_time, Computer, user, source_type, process_name, task_detail, risk
| rename task_time as "Time", Computer as "Host", user as "User", source_type as "Source", process_name as "Process", task_detail as "Task Details", risk as "Risk"
| sort - Time</query>
          <earliest>$time_range.earliest$</earliest>
          <latest>$time_range.latest$</latest>
        </search>
        <option name="count">20</option>
        <option name="drilldown">row</option>
        <option name="wrap">true</option>
        <format type="color" field="Risk">
          <colorPalette type="map">{"HIGH":#DC4E41,"MEDIUM":#F8BE34}</colorPalette>
        </format>
      </table>
    </panel>
  </row>

  <!-- ===================== ROW 4: New User Accounts Created (Event 4720) ===================== -->
  <row>
    <panel>
      <title>New User Accounts Created (Security Event 4720)</title>
      <table>
        <search>
          <query>index=wineventlog EventCode=4720
| eval creation_time=strftime(_time, "%Y-%m-%d %H:%M:%S")
| eval risk=case(
    match(TargetUserName, "(?i)(admin|svc|service|test|temp|guest)"), "HIGH",
    SubjectUserName != "SYSTEM" AND NOT match(SubjectUserName, "(?i)admin"), "HIGH",
    1=1, "MEDIUM"
    )
| table creation_time, ComputerName, SubjectUserName, SubjectDomainName, TargetUserName, TargetDomainName, risk
| rename creation_time as "Time", ComputerName as "Host", SubjectUserName as "Created By", SubjectDomainName as "Creator Domain", TargetUserName as "New Account", TargetDomainName as "Account Domain", risk as "Risk"
| sort - Time</query>
          <earliest>$time_range.earliest$</earliest>
          <latest>$time_range.latest$</latest>
        </search>
        <option name="count">20</option>
        <option name="drilldown">row</option>
        <format type="color" field="Risk">
          <colorPalette type="map">{"HIGH":#DC4E41,"MEDIUM":#F8BE34}</colorPalette>
        </format>
      </table>
    </panel>
  </row>

  <!-- ===================== ROW 5: Group Membership Changes ===================== -->
  <row>
    <panel>
      <title>Group Membership Changes (Security Events 4728 / 4732 / 4756)</title>
      <table>
        <search>
          <query>index=wineventlog EventCode IN (4728, 4732, 4756)
| eval change_time=strftime(_time, "%Y-%m-%d %H:%M:%S")
| eval group_scope=case(
    EventCode=4728, "Global Security Group",
    EventCode=4732, "Local Security Group",
    EventCode=4756, "Universal Security Group"
    )
| eval risk=case(
    match(TargetUserName, "(?i)(domain admins|enterprise admins|schema admins|administrators|backup operators|account operators)"), "CRITICAL",
    match(TargetUserName, "(?i)(remote desktop|dnsadmins|group policy|server operators)"), "HIGH",
    1=1, "MEDIUM"
    )
| table change_time, ComputerName, SubjectUserName, MemberName, TargetUserName, group_scope, risk
| rename change_time as "Time", ComputerName as "Host", SubjectUserName as "Changed By", MemberName as "Member Added", TargetUserName as "Group Name", group_scope as "Group Scope", risk as "Risk"
| sort - Time</query>
          <earliest>$time_range.earliest$</earliest>
          <latest>$time_range.latest$</latest>
        </search>
        <option name="count">20</option>
        <option name="drilldown">row</option>
        <format type="color" field="Risk">
          <colorPalette type="map">{"CRITICAL":#DC4E41,"HIGH":#F1813F,"MEDIUM":#F8BE34}</colorPalette>
        </format>
      </table>
    </panel>
  </row>
</dashboard>
