# sigma/privilege-escalation/scheduled-task-created.yml
title: Suspicious Scheduled Task Creation via schtasks.exe
id: e1b3f8a2-9d4c-4e6f-f7a8-1c2b3d4e5f6a
status: experimental
level: medium
description: |
    Detects the creation of scheduled tasks via schtasks.exe with suspicious
    characteristics. Attackers use scheduled tasks for persistence (surviving
    reboots), privilege escalation (running as SYSTEM), and execution
    (triggering remotely). This rule flags tasks running from temp directories,
    using script interpreters, executing as SYSTEM, or created via remote
    execution frameworks.
references:
    - https://attack.mitre.org/techniques/T1053/005/
    - https://docs.microsoft.com/en-us/windows/win32/taskschd/schtasks
author: Kud Lee
date: 2026/03/19
tags:
    - attack.privilege_escalation
    - attack.persistence
    - attack.execution
    - attack.t1053.005
logsource:
    category: process_creation
    product: windows
detection:
    selection:
        EventID: 1
        Image|endswith: '\schtasks.exe'
        CommandLine|contains: '/create'
    suspicious_system:
        CommandLine|contains:
            - '/ru system'
            - '/ru SYSTEM'
            - '/ru "SYSTEM"'
    suspicious_path:
        CommandLine|contains:
            - '\Temp\'
            - '\tmp\'
            - '\AppData\'
            - '\Public\'
            - '\Downloads\'
    suspicious_command:
        CommandLine|contains:
            - 'powershell'
            - 'cmd.exe /c'
            - 'mshta'
            - 'certutil'
            - 'bitsadmin'
            - 'regsvr32'
    suspicious_script:
        CommandLine|contains:
            - '.bat'
            - '.ps1'
            - '.vbs'
            - '.js'
            - '.hta'
    suspicious_parent:
        ParentImage|endswith:
            - '\wmiprvse.exe'
            - '\PSEXESVC.exe'
            - '\wsmprovhost.exe'
    condition: selection and (suspicious_system or suspicious_path or suspicious_command or suspicious_script or suspicious_parent)
falsepositives:
    - Legitimate system administration task scheduling
    - Software installations creating maintenance tasks
    - Group Policy deploying scheduled tasks
    - SCCM and Intune management task deployments
